# apps/runtime/config/java/supervisord.conf
[supervisord]
nodaemon=false

[program:tomcat]
command=xxx
autostart=false
directory=/usr/local/tomcat
