# apps/runtime/config/java/supervisord.conf
[supervisord]
nodaemon=false

[program:jenkins]
command=java -jar jenkins.war --httpPort=8080
autostart=false
directory=/data/apps/jenkins
