# apps/runtime/config/java/supervisord.conf
[supervisord]
nodaemon=false

[program:jenkins]
command=java -jar jenkins.war
autostart=false
directory=/data/apps/jenkins
