# apps/runtime/config/java/supervisord.conf
[supervisord]
nodaemon=false

[program:tomcat]
command=catalina.sh run
autostart=false
directory=/data/apps/tomcat/bin
